# linux/debian/etc/cron.weekly/podman-prune
#!/bin/sh

name=$(basename "$0")
logger "$name: pruning dangling images and volumes"

podman image prune -f
podman volume prune -f

logger "$name: done!"
